# tomcat/supervisord.conf
[supervisord]
nodaemon=true

[program:http]
command=/opt/tomcat/bin/catalina.sh start
numprocs=1
autostart=true
autorestart=true
